# deploy/tailscale/cert-renew.timer
[Unit]
Description=Monthly renewal of Tailscale cert for clideck

[Timer]
# Tailscale certs are 90 days; monthly renewal keeps well ahead.
OnCalendar=monthly
RandomizedDelaySec=6h
Persistent=true
Unit=clideck-cert-renew.service

[Install]
WantedBy=timers.target
