# deploy/tailscale/clideck.service
[Unit]
Description=clideck (Tailscale HTTPS)
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
Type=simple

# Run as a dedicated unprivileged user. setup.sh creates this user with
# `useradd --system --create-home --home /var/lib/clideck --shell /usr/sbin/nologin clideck`.
# Running as root would silently negate ProtectSystem / ProtectHome / Cap*.
User=clideck
Group=clideck

# Stable HOME so clideck's per-user state lands under /var/lib/clideck/.clideck
Environment=HOME=/var/lib/clideck

# Env file is written by deploy/tailscale/setup.sh. Contents:
#   CLIDECK_HOST=<tailnet-ip>
#   CLIDECK_PORT=443
#   CLIDECK_TLS_CERT_PATH=/etc/clideck/tls/fullchain.pem
#   CLIDECK_TLS_KEY_PATH=/etc/clideck/tls/privkey.pem
#   CLIDECK_ALLOWED_ORIGINS="https://<host>.<tailnet>.ts.net"
#   CLIDECK_TAILSCALE_FQDN=<host>.<tailnet>.ts.net
#   CLIDECK_PUBLIC_MODE=1
#   CLIDECK_USERNAME="<your-user>"
#   CLIDECK_PASSWORD="<strong-password>"
#   CLIDECK_COOKIE_SECURE=1
EnvironmentFile=/etc/clideck/clideck.env

ExecStart=/usr/bin/env clideck

# Graceful shutdown: clideck persists sessions on SIGTERM.
KillSignal=SIGTERM
TimeoutStopSec=15
Restart=on-failure
RestartSec=3

# Hardening — these are now meaningful because User= is set.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
PrivateTmp=true
PrivateDevices=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native

# Only the directories we genuinely write to.
ReadWritePaths=/etc/clideck /var/lib/clideck

# Required to bind :443 as the unprivileged clideck user.
# setup.sh also runs `setcap cap_net_bind_service=+ep` on the node binary;
# this directive grants the *capability* to the process, while setcap
# permits the binary to keep it after exec().
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
